# flaskAutomation/scripts/flask_project.conf
[program:flask_project]
command = gunicorn app:app -b 0.0.0.0:8000
directory = /home/www/flask_project
user = newuser